# client-builder/src/modules/mindbeamer/filesystem/lib/systemd/system/mindbeamer_provision.service
[Unit]
Description=Start Mindbeamer provision script

[Service]
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/bin/bash /opt/mindbeamer/provision.sh
Restart=always
RestartSec=10s
KillMode=process
TimeoutSec=infinity

[Install]
WantedBy=graphical.target