# client-builder/src/modules/mindbeamer/filesystem/etc/cron.d/mindbeamer
# /etc/cron.d/mindbeamer: crontab fragment for mindbeamer

@reboot root /bin/bash /opt/mindbeamer/wpa-supplicant.sh
@reboot root /bin/bash /opt/mindbeamer/node-splash.sh
*/5 * * * * root /bin/bash /opt/mindbeamer/node-splash.sh
#* * * * * pi sleep 5 && /bin/bash /home/pi/scripts/refresh
@reboot root sleep 120 && /bin/bash /opt/mindbeamer/provision.sh
@reboot root timedatectl set-timezone America/New_York
